# systemd/i3-resurrect-save.timer
[Unit]
Description=i3-resurrect save all workspaces (timer)

[Timer]
OnUnitActiveSec=60s
OnBootSec=60s

[Install]
WantedBy=timers.target
